# examples/mobilizon-bot.timer
[Unit]
Requires=mobilizon-bot.service

[Timer]
Unit=mobilizon-bot.service
OnCalendar=hourly

[Install]
WantedBy=timers.target

